# files/system/etc/systemd/system/disable-usb-wakeups.service
[Unit]
Description=Disable Spurious USB and PCIe Wakeup Events
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'if [ -f /proc/acpi/wakeup ]; then for dev in XHC XHCI TXHC TBT0 TBT1 GLAN; do grep -qE "^${dev}\b.*enabled" /proc/acpi/wakeup && echo "$dev" > /proc/acpi/wakeup || true; done; fi'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
