# files/system/etc/systemd/system/ultra-hardware-init.service
[Unit]
Description=Core Ultra Hardware Initialization and Tuning Engine
Before=sysinit.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/ultra-hardware-init.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
