# files/system/etc/systemd/system/ultra-hardware-init.service
[Unit]
Description=Core Ultra Hardware Initialization and Tuning Engine
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ultra-hardware-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
