# files/system/etc/profile.d/gpu-acceleration.sh
#!/bin/sh
export LIBVA_DRIVER_NAME=iHD
export MESA_SHADER_CACHE_MAX_SIZE=8G
export MESA_DISK_CACHE_SINGLE_FILE=1
