# files/system/etc/udev/rules.d/99-scx-power.rules
# Trigger the scheduler switch script immediately when the power cable is pulled or connected
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/usr/local/bin/scx-power-switch.sh"