# files/system/etc/systemd/user/tracker-miner-rss-3.service
[Unit]
Description=Disabled File System Indexing Daemon
[Service]
Type=oneshot
ExecStart=/bin/true